# opinions_scrape.cron
0 8/5 * * * cd /home/rahmi/Documents/newsNLUAPP/ && ./opinions_scrape.py
